# bacchus-lab-env/debian/postinst
#!/bin/sh

set -e

systemctl enable gdm.service
systemctl disable NetworkManager.service

exit 0
