# bacchus-shutdown/debian/postinst
#!/bin/bash

set -e

systemctl enable bacchus-shutdown.timer

exit 0
