# bacchus-exam-env/debian/bacchus-exam-env.postinst
#!/bin/sh

systemctl enable --now bacchus-exam-mode-enable.timer bacchus-exam-mode-disable.timer
